# systemd/dnsmasq.service.d/jog-pre.conf
[Unit]
Requires=jog-render-dnsmasq.service
After=jog-render-dnsmasq.service

[Service]
ExecStartPre=/usr/bin/systemctl start jog-render-dnsmasq.service
